# wifi-reconnect.service
[Unit]
Description=WiFi Reconnect Watchdog
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
EnvironmentFile=/etc/default/wifi-reconnect
ExecStart=/usr/local/bin/wifi-reconnect.sh $TARGET_SSID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
